# deb/install/prerm
#!/bin/bash

systemctl stop waggle-network-watchdog.service || true
systemctl disable waggle-network-watchdog.service || true
