# deb/install/prerm
#!/bin/bash

systemctl stop waggle-network-watchdog.timer || true
systemctl disable waggle-network-watchdog.timer || true
